# roles/configure/templates/leihs_create-master-secret.service
[Unit]
Description=Leihs create-master-secret

[Service]
Type=oneshot
# RemainAfterExit=yes
ExecStart=/leihs/deploy/bin/create-master-secret
