# roles/configure/templates/leihs_create-master-secret.path
[Path]
PathChanged=/leihs/config/
PathChanged=/leihs/config/master-secret.txt
PathModified=/leihs/config/master-secret.txt
PathExists=/leihs/config/master-secret.txt
