# roles/configure/templates/leihs_set-secret-key-base-for-leihs-legacy.path
[Path]
PathModified=/leihs/config/master-secret.txt
PathExists=/leihs/config/master-secret.txt
