# roles/db-sweep-service/templates/leihs-db-sweep.service
[Unit]
Description=Leihs Database Sweep


[Service]
SyslogIdentifier=leihs-db-sweep
Type=oneshot
ExecStart=/usr/local/bin/leihs-db-sweep
