# roles/db-backup/templates/leihs-db-backup.timer
[Unit]
Description=Leihs Daily Backup Timer


[Timer]
OnCalendar=*-*-* 04:00:00
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
WantedBy=leihs-migration.service
